# supervisor-app.conf
[supervisord]
[program:app-uwsgi]
command = /usr/bin/uwsgi --ini /home/docker/code/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx

[program:daphne]
directory=/home/docker/code/news
command=daphne -b 127.0.0.1 -p 8002 --proxy-headers news.asgi:application
autostart=true
autorestart=true
stdout_logfile=/tmp/websocket.log  #Journal
redirect_stderr=true
